# docker/build/supervisord.conf
[supervisord]
nodaemon=true

[program:ant]
directory="$APP_PATH"
command=/bin/bash -c "gosu "$HOST_USER" ant -emacs work"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
autorestart=false

[program:apache2]
command=/bin/bash -c "/sbin/my_init"

[program:watch]
command=/bin/bash -c "gosu "$HOST_USER" /scripts/watch-app.sh"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
