# docker/build/supervisord.conf
[supervisord]
nodaemon=true

[program:ant]
directory=/var/www/html
command=/bin/bash -c "gosu "$HOST_USER" ant -emacs work"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=false

[program:apache2]
command=/bin/bash -c "apache2-foreground"
